# deploy/photonicat/context-pilot.init
#!/bin/sh /etc/rc.common
# Context Pilot orchestrator — OpenWrt procd service (Photonicat, no Docker).
#
# Install as /etc/init.d/context-pilot, then:  /etc/init.d/context-pilot enable
# Logs go to the system log: `logread -e context-pilot`.

USE_PROCD=1
# Start late (after network/dnsmasq); stop early.
START=95
STOP=10

PROG=/mnt/data/context-pilot/bin/cp-orchestrator
CP_HOME=/mnt/data/context-pilot/home
CP_WEB=/mnt/data/context-pilot/web
# Agent (cp TUI) binary the supervisor is allowed to spawn for dashboard-created
# agents. Spawned agents inherit this process's env (HOME, provider keys, …).
CP_AGENT=/mnt/data/context-pilot/bin/tui

start_service() {
    # Persisted state (registry, auth.db, central config, Claude creds, agent
    # project roots) all live under $HOME on the big data partition.
    mkdir -p "$CP_HOME"

    # Provisioning values (templated by Ansible, chmod 600), sourced at boot:
    #   seed.env      — per-UNIT admin seed (CP_SEED_ADMIN_*); write-once. The
    #                   orchestrator seeds the first Admin only when the user
    #                   table is empty (idempotent).
    #   providers.env — per-CLIENT provider/integration API keys; rewritten on
    #                   each Ansible run from the client's vault. Placed in the
    #                   process env so agents can use the model and the cockpit
    #                   Services pane lights up. Two clients → two key sets.
    CP_SEED_ENV=/mnt/data/context-pilot/seed.env
    [ -f "$CP_SEED_ENV" ] && . "$CP_SEED_ENV"
    CP_PROVIDERS_ENV=/mnt/data/context-pilot/providers.env
    [ -f "$CP_PROVIDERS_ENV" ] && . "$CP_PROVIDERS_ENV"

    # IT maintenance plane. Caddy terminates TLS and fronts the public :9090
    # (private-CA `tls internal`); the orchestrator's own maintenance listener
    # therefore binds a LOOPBACK-only internal port (9191) that only Caddy
    # reaches — it is never exposed directly. The application LAN-only guard
    # (CP_MAINT_LAN_ONLY, on by default) is a second layer should the bind ever
    # widen. The paper-admin window is mitigated by all three layers (loopback/
    # LAN-only + unit-unique password + forced change), per local-tls-onboarding.
    CP_MAINT_PORT=9191
    CP_MAINT_BIND=127.0.0.1

    # Caddy integration (M3): the orchestrator regenerates this Caddyfile (name/
    # IP, provisioned gate) and reloads Caddy through its admin API. Point both
    # at the on-box Caddy so identity/finalize can re-issue the leaf and start
    # serving the cockpit on :443.
    CP_CADDYFILE=/mnt/data/context-pilot/caddy/Caddyfile
    CP_CADDY_BIN=/mnt/data/context-pilot/bin/caddy

    procd_open_instance
    procd_set_param command "$PROG"
    procd_set_param env \
        HOME="$CP_HOME" \
        CP_AUTH_ENABLED=1 \
        CP_ORCH_PORT=7878 \
        CP_MAINT_PORT="${CP_MAINT_PORT:-9191}" \
        CP_MAINT_BIND="${CP_MAINT_BIND:-127.0.0.1}" \
        CP_CADDYFILE="$CP_CADDYFILE" \
        CP_CADDY_BIN="$CP_CADDY_BIN" \
        CP_WEB_ROOT="$CP_WEB" \
        CP_AGENT_BINARY="$CP_AGENT" \
        CP_SEED_ADMIN_EMAIL="${CP_SEED_ADMIN_EMAIL:-}" \
        CP_SEED_ADMIN_NAME="${CP_SEED_ADMIN_NAME:-}" \
        CP_SEED_ADMIN_PASSWORD="${CP_SEED_ADMIN_PASSWORD:-}" \
        CP_SEED_ADMIN_PASSWORD_FILE="${CP_SEED_ADMIN_PASSWORD_FILE:-}" \
        ANTHROPIC_API_KEY="${ANTHROPIC_API_KEY:-}" \
        DEEPSEEK_API_KEY="${DEEPSEEK_API_KEY:-}" \
        GROQ_API_KEY="${GROQ_API_KEY:-}" \
        XAI_API_KEY="${XAI_API_KEY:-}" \
        MINIMAX_API_KEY="${MINIMAX_API_KEY:-}" \
        VOYAGE_API_KEY="${VOYAGE_API_KEY:-}" \
        DATALAB_API_KEY="${DATALAB_API_KEY:-}" \
        BRAVE_API_KEY="${BRAVE_API_KEY:-}" \
        FIRECRAWL_API_KEY="${FIRECRAWL_API_KEY:-}" \
        GITHUB_TOKEN="${GITHUB_TOKEN:-}"
    # Restart on crash, but back off if it crash-loops.
    procd_set_param respawn 3600 5 5
    procd_set_param stdout 1
    procd_set_param stderr 1
    procd_close_instance
}
